# wayland/.config/systemd/user/hyprpaper.service
[Unit]
Description=Hyprpaper

[Service]
Type=simple
ExecStart=/usr/bin/hyprpaper
Restart=always
TimeoutStartSec=5

[Install]
WantedBy=default.target
