# system/udev-rules.d/99-system76.rules
ACTION=="add", SUBSYSTEM=="leds", KERNEL=="system76_acpi::kbd_backlight", RUN+="/bin/chmod a+w /sys/class/leds/%k/brightness"


